# files/service/hazelcast/run
#!/bin/bash

exec /sbin/setuser pentaho /opt/hazelcast/launch-hazelcast.sh
